# etc/profile.d/omd.sh
#!/bin/bash

export OPENMD_UTILS_MACRO_DIR="/usr/share/openmd-utils/macros"
